# examples/systemd/reefer-stop.target
[Unit]
Description=Stop reefer service

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop reefer.service
