# examples/systemd/reefer-start.timer
[Unit]
Description=Start reefer service every morning

[Timer]
OnCalendar=06:00:00
Unit=reefer-start.target
