# examples/systemd/reefer-mailer.service
[Unit]
Description=Send reef status email

[Service]
Type=oneshot
ExecStart=/usr/bin/reefer-mailer
