# examples/systemd/reefer-start.target
[Unit]
Description=Start reefer service

[Service]
Type=oneshot
ExecStart=/bin/systemctl start reefer.service
